# res/grommunio-certbot-renew.timer
[Unit]
Description=Weekly certbot renewal check

[Timer]
OnCalendar=weekly
AccuracySec=1h
Persistent=true

[Install]
WantedBy=timers.target
