# res/grommunio-fetchmail.timer
[Unit]
Description=fetchmail configuration generation (hourly)

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
